# contrib/init/cpuchaind.conf
description "CPUchain Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env CPUCHAIND_BIN="/usr/bin/cpuchaind"
env CPUCHAIND_USER="cpuchain"
env CPUCHAIND_GROUP="cpuchain"
env CPUCHAIND_PIDDIR="/var/run/cpuchaind"
# upstart can't handle variables constructed with other variables
env CPUCHAIND_PIDFILE="/var/run/cpuchaind/cpuchaind.pid"
env CPUCHAIND_CONFIGFILE="/etc/cpuchain/cpuchain.conf"
env CPUCHAIND_DATADIR="/var/lib/cpuchaind"

expect fork

respawn
respawn limit 5 120
kill timeout 600

pre-start script
    # this will catch non-existent config files
    # cpuchaind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$CPUCHAIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run cpuchaind."
        echo "The setting must appear in $CPUCHAIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Bitcoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$CPUCHAIND_PIDDIR"
    chmod 0755 "$CPUCHAIND_PIDDIR"
    chown $CPUCHAIND_USER:$CPUCHAIND_GROUP "$CPUCHAIND_PIDDIR"
    chown $CPUCHAIND_USER:$CPUCHAIND_GROUP "$CPUCHAIND_CONFIGFILE"
    chmod 0660 "$CPUCHAIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$CPUCHAIND_PIDFILE" \
    --chuid $CPUCHAIND_USER:$CPUCHAIND_GROUP \
    --exec "$CPUCHAIND_BIN" \
    -- \
    -pid="$CPUCHAIND_PIDFILE" \
    -conf="$CPUCHAIND_CONFIGFILE" \
    -datadir="$CPUCHAIND_DATADIR" \
    -disablewallet \
    -daemon

